# templates/zookeeper.service.j2
[Unit]
Description=Zookeeper
Before=
After=network.target

[Service]
Type=forking
User={{ zookeeper_user }}
Group={{ zookeeper_group }}
Environment=ZOO_LOG_DIR={{ zookeeper_log_dir }}
ExecStart={{ zookeeper_install_dir }}/bin/zkServer.sh start
ExecStop={{ zookeeper_install_dir }}/bin/zkServer.sh stop
ExecReload={{ zookeeper_install_dir }}/bin/zkServer.sh restart

[Install]
WantedBy=multi-user.target